# rpb-debian-iot/conf/bluetooth_6lowpand.service
[Unit]
Description=Bluetooth 6lowpan joiner service
After=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bluetooth_6lowpand
